# share/systemd/yca-ocsp.service
# systemd unit for yca-ocsp, the delegated OCSP responder.
#
# The responder must run as the store owner: the DB is 0600, reading a WAL
# database needs write access to its -wal/-shm siblings.
#
# Credential rotation (yca create ocsp) needs no restart - the responder
# reloads the cert/key when the files change.

[Unit]
Description=yca OCSP responder (delegated, RFC 6960)
After=network-online.target
Wants=network-online.target

[Service]
Type=exec
User=yca
Group=yca
ExecStart=/usr/local/bin/yca-ocsp --store /var/lib/yca/store --listen 127.0.0.1:23509
Restart=on-failure
RestartSec=2s

# State dir: rw; network: AF_INET/AF_INET6 (listener)
UMask=0077
NoNewPrivileges=yes
ProtectSystem=strict
ReadWritePaths=/var/lib/yca/store
ProtectHome=yes
PrivateTmp=yes
PrivateDevices=yes
ProtectClock=yes
ProtectControlGroups=yes
ProtectKernelLogs=yes
ProtectKernelModules=yes
ProtectKernelTunables=yes
ProtectProc=invisible
RestrictAddressFamilies=AF_INET AF_INET6
RestrictNamespaces=yes
RestrictRealtime=yes
RestrictSUIDSGID=yes
LockPersonality=yes
MemoryDenyWriteExecute=yes
CapabilityBoundingSet=
SystemCallArchitectures=native
SystemCallFilter=@system-service
SystemCallErrorNumber=EPERM

[Install]
WantedBy=multi-user.target
