# share/systemd/yca-ocsp-rotate.timer
# systemd timer driving yca-ocsp-rotate.service.
#
# Cadence vs validity: the responder certificate lives ocsp_valid_days
# (default 7). Twice a week leaves 3-4 days of margin; keep the interval
# comfortably below the validity if you change either.

[Unit]
Description=Periodic rotation of the yca OCSP responder certificate

[Timer]
OnCalendar=Mon,Thu 03:00
RandomizedDelaySec=1h
Persistent=true

[Install]
WantedBy=timers.target
