# share/systemd/yca-crl-refresh.timer
# systemd timer driving yca-crl-refresh.service (signing CRL only).
#
# Cadence vs nextUpdate: the signing CRL promises re-publication within
# crl_next_update_days (7). Daily leaves 6 days of margin; keep the
# interval comfortably below the horizon if you change either. The root
# CRL has its own timer - see yca-root-crl-refresh.timer.

[Unit]
Description=Periodic re-signing of the published yca signing CRL

[Timer]
OnCalendar=*-*-* 03:00
RandomizedDelaySec=1h
Persistent=true

[Install]
WantedBy=timers.target
