# share/systemd/yca-crl-refresh.service
# systemd service re-signing the published yca signing CRL.
# Started by yca-crl-refresh.timer.
#
# `yca refresh crl signing` re-signs the signing CRL from its current
# entry set: same revocations, crlNumber+1, fresh thisUpdate/nextUpdate.
# Without it a strict CRL checker considers the signing CRL stale
# crl_next_update_days (7) after the last revocation. The root CRL runs
# on its own, much slower cadence - see yca-root-crl-refresh.{service,timer}.

[Unit]
Description=Re-sign the published yca signing CRL

[Service]
Type=oneshot
User=yca
Group=yca
EnvironmentFile=/etc/yca/yca.env
ExecStart=/usr/local/bin/yca --config /var/lib/yca/yca.toml --store /var/lib/yca/store refresh crl signing

# State dir: rw; pkcs11 backend: AF_UNIX (pcscd); no network.
UMask=0077
NoNewPrivileges=yes
ProtectSystem=strict
ReadWritePaths=/var/lib/yca
ProtectHome=yes
PrivateTmp=yes
PrivateDevices=yes
ProtectClock=yes
ProtectControlGroups=yes
ProtectKernelLogs=yes
ProtectKernelModules=yes
ProtectKernelTunables=yes
ProtectProc=invisible
RestrictAddressFamilies=AF_UNIX
RestrictNamespaces=yes
RestrictRealtime=yes
RestrictSUIDSGID=yes
LockPersonality=yes
MemoryDenyWriteExecute=yes
CapabilityBoundingSet=
SystemCallArchitectures=native
SystemCallFilter=@system-service
SystemCallErrorNumber=EPERM
